# Tools/systemd/pai-backup.service
[Unit]
Description=PAI Weekly Backup
After=network.target

[Service]
Type=oneshot
ExecStart=/home/user/Personal_AI_Infrastructure/Tools/backup-pai.sh --full
User=%i
StandardOutput=append:/home/user/pai-backups/backup.log
StandardError=append:/home/user/pai-backups/backup.log

[Install]
WantedBy=multi-user.target
