# Tools/systemd/pai-backup.timer
[Unit]
Description=Run PAI Backup Weekly

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
